# ci/ansible/roles/docker_registry/templates/docker-registry.mount.j2
[Unit]
Description=k3s docker-registry data files
Before=docker-registry.service

[Mount]
What={{ datadir }}/docker-registry
Where=/var/lib/docker-registry
Type=none
Options=bind
